# deploy/vps-sentry.service
[Unit]
Description=VPS Sentry — resource monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vps-sentry
EnvironmentFile=/opt/vps-sentry/.env
Environment=PYTHONPATH=/opt/vps-sentry/src
ExecStart=/opt/vps-sentry/.venv/bin/python -m vps_sentry --config /opt/vps-sentry/config.yml
Restart=on-failure
RestartSec=10s
User=vps-sentry
Nice=10
MemoryMax=64M

[Install]
WantedBy=multi-user.target
